# docker/supervisord.ini
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:proxy]
command=node /app/proxy/dist/index.js
directory=/app/proxy
autostart=true
autorestart=true
environment=PROXY_PORT=3001,DB_PATH=%(ENV_DB_PATH)s
startsecs=3
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
